# scripts/scrape.cron
0 8,16 * * * utube-scraper.py >> /home/subrat/.local/utube/cronjob-hist.txt
