# scripts/scrape.cron
0 8,16 * * * utube-scraper >> /home/subrat/.local/utube/cronjob-hist.txt
